# src/s6/debian-root/etc/s6-overlay/s6-rc.d/stubby/finish
#!/usr/bin/env sh

s6-echo "Stopping stubby"
killall -9 stubby
s6-echo "Stopped stubby"
