# src/s6/debian-root/etc/s6-overlay/s6-rc.d/stubby/run
#!/usr/bin/env sh

s6-echo "Starting stubby"
cp --no-clobber /tmp/config/stubby.yml /config/stubby.yml
stubby -C /config/stubby.yml
