# src/s6/debian-root/etc/s6-overlay/s6-rc.d/cloudflared/finish
#!/usr/bin/env sh

s6-echo "Stopping cloudflared"
killall -9 cloudflared
s6-echo "Stopped cloudflared"
