# src/s6/debian-root/etc/s6-overlay/s6-rc.d/cloudflared/run
#!/usr/bin/env sh

s6-echo "Starting cloudflared"

cp --no-clobber /tmp/config/cloudflared.yml /config/cloudflared.yml
cloudflared --config /config/cloudflared.yml
